# etc-init/digidj-server-web.conf
start on starting digidj-server
stop on stopping digidj-server
